# FV3GFSwfm/rt_hr3-chem/rt_hr3-chem.crontab

#################### rt_hr3-chem ####################
MAILTO=""
*/5 * * * * /apps/rocoto/1.3.7/bin/rocotorun -d /scratch1/BMC/gsd-fv3/rtruns/UFS-Chem/FV3GFSwfm/rt_hr3-chem/rt_hr3-chem.db -w /scratch1/BMC/gsd-fv3/rtruns/UFS-Chem/FV3GFSwfm/rt_hr3-chem/rt_hr3-chem.xml
#################################################################
